# vvv.conf
[program:vvv]
command=python3 /home/curator2/vvv/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/vvv/test.err.log
stdout_logfile=/home/curator2/vvv/test.out.log